# etc/s6-overlay/s6-rc.d/03-avahi/run
#!/command/with-contenv sh
echo "Starting avahi"
exec s6-notifyoncheck avahi-daemon --no-chroot
