# etc/s6-overlay/s6-rc.d/04-nqptp/run
#!/command/with-contenv sh
echo "Starting nqptp"
exec /usr/local/bin/nqptp
